# contrib/init/cryptomilesd.conf
description "Cryptomiles Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CRYPTOMILESD_BIN="/usr/bin/cryptomilesd"
env CRYPTOMILESD_USER="cryptomiles"
env CRYPTOMILESD_GROUP="cryptomiles"
env CRYPTOMILESD_PIDDIR="/var/run/cryptomilesd"
# upstart can't handle variables constructed with other variables
env CRYPTOMILESD_PIDFILE="/var/run/cryptomilesd/cryptomilesd.pid"
env CRYPTOMILESD_CONFIGFILE="/etc/cryptomiles/cryptomiles.conf"
env CRYPTOMILESD_DATADIR="/var/lib/cryptomilesd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # cryptomilesd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CRYPTOMILESD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run cryptomilesd."
        echo "The setting must appear in $CRYPTOMILESD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Cryptomiles Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CRYPTOMILESD_PIDDIR"
    chmod 0755 "$CRYPTOMILESD_PIDDIR"
    chown $CRYPTOMILESD_USER:$CRYPTOMILESD_GROUP "$CRYPTOMILESD_PIDDIR"
    chown $CRYPTOMILESD_USER:$CRYPTOMILESD_GROUP "$CRYPTOMILESD_CONFIGFILE"
    chmod 0660 "$CRYPTOMILESD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CRYPTOMILESD_PIDFILE" \
    --chuid $CRYPTOMILESD_USER:$CRYPTOMILESD_GROUP \
    --exec "$CRYPTOMILESD_BIN" \
    -- \
    -pid="$CRYPTOMILESD_PIDFILE" \
    -conf="$CRYPTOMILESD_CONFIGFILE" \
    -datadir="$CRYPTOMILESD_DATADIR" \
    -disablewallet \
    -daemon

